# images/php/crontab
* * * * * root cd /var/www/beautybox-php && php artisan schedule:run >> /dev/null 2>&1
